# jetson/webservices-port.service
[Unit]
Description=farm-ng web services socket

[Socket]
ListenStream=80
NoDelay=true
